# tunnels/rootfs/etc/services.d/tunnels/run
#!/usr/bin/with-contenv bashio

declare -a options
KEYSTORE="/data/ssh-keystore"

options+=(-T)
options+=(-p $(bashio::config 'port'))
options+=(-i "$KEYSTORE/tunnels_rsa")
options+=(-o StrictHostKeyChecking=no)
options+=(-o ServerAliveInterval=$(bashio::config 'ServerKeepAliveInterval'))

for id in $(bashio::config "tunnels|keys"); do
  tunnel=$(bashio::config "tunnels[${id}].tunnel")
  bashio::log.info "Adding tunnel to request: $tunnel"
  options+=(-R $tunnel)
done

SERVER=$(bashio::config 'server')

USER=$(bashio::config 'ssh_user')
options+=($USER@$SERVER)

if [[ $(bashio::config.true 'sni_proxy') ]]; then
  options+=(sni-proxy=true)
fi

if [[ $(bashio::config.true 'tcp_alias') ]]; then
  options+=(tcp-alias=true)
fi

bashio::log.info "Command to execute:"
bashio::log.info "ssh ${options[@]}"

bashio::log.info "Attempting to establish SSH tunnels..."
exec ssh ${options[@]}
